# supervisor-cron.conf
[program:cron]
autostart = true
autorestart = true
command=cron -f
stderr_logfile = /var/log/cron-stderr.log
stdout_logfile = /var/log/cron-stdout.log
startretries = 3
